# private_Library/LaunchAgents/com.micahsvenson.obsidian-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.micahsvenson.obsidian-sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>waited=0; while ! brctl status com.apple.CloudDocs 2>/dev/null | grep -q "caught-up"; do sleep 5; waited=$((waited + 5)); if [ $waited -ge 60 ]; then echo "iCloud still syncing after 60s, skipping"; exit 0; fi; done; cd ~/Projects/Obsidian || exit 0; git fetch origin main 2>/dev/null || exit 0; git merge --ff-only origin/main 2>/dev/null || true; git add -A; if ! git diff --cached --quiet; then git commit -m "auto: sync $(date +%Y-%m-%d_%H:%M)"; git push origin main; fi</string>
    </array>
    <key>StartInterval</key>
    <integer>300</integer>
    <key>StandardOutPath</key>
    <string>/Users/micahsvenson/Library/Logs/obsidian-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/micahsvenson/Library/Logs/obsidian-sync.log</string>
    <key>RunAtLoad</key>
    <true/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/micahsvenson</string>
    </dict>
</dict>
</plist>
